# DEBIAN/postinst
#!/bin/bash

path=`dirname $0`

# Fonction: installation
installation() {
<<<<<<< HEAD
  cp /etc/default/openplacos /etc/default/openplacos.back
  cp /usr/lib/ruby/openplacos/server/config_with_VirtualPlacos.yaml /etc/default/openplacos

  # Pre-requis
  echo "----------------------------------------------------"
  echo "Ruby gem lib installation "
  echo "This could take about 10 min -- please wait"
  echo "----------------------------------------------------"
  gem install activerecord mysql serialport openplacos ruby-dbus-openplacos micro-optparse --no-ri --no-rdoc
=======
    # Pre-requis
    echo "----------------------------------------------------"
    echo "Ruby gem lib installation "
    echo "This could take about 10 min -- please wait"
    echo "----------------------------------------------------"
    gem install activerecord mysql serialport openplacos micro-optparse choice rails --bindir /usr/bin --no-ri --no-rdoc
>>>>>>> unstable
    
  # User openplacos
    adduser openplacos --system -disabled-login -no-create-home


  # link into path
    ln -s -f /usr/lib/ruby/openplacos/server/main.rb /usr/bin/openplacos-server
    ln -s -f /usr/lib/ruby/openplacos/client/CLI_client/opos-client.rb /usr/bin/openplacos
    ln -s -f /usr/lib/ruby/openplacos/client/deprecated/gtk/gtk.rb /usr/bin/openplacos-gtk
    ln -s -f /usr/lib/ruby/openplacos/client/xml-rpc/xml-rpc-client.rb  /usr/bin/openplacos-xmlrpc
    ln -s -f /usr/lib/ruby/openplacos/client/soap/soap-client.rb  /usr/bin/openplacos-soap

  # dbus integration
    cp -f /usr/lib/ruby/openplacos/setup_files/openplacos.conf /etc/dbus-1/system.d/

  # reload dbus
    /etc/init.d/dbus reload

  # Start with system
    cp -f /usr/lib/ruby/openplacos/setup_files/openplacos /etc/init.d/
    update-rc.d openplacos defaults 98 02

}

if [ "$1" = "configure" ]; then
    installation
fi
